# config/rc.local
#!/bin/bash
if [ -f "/root/postinstall.sh" ]; then
    bash /root/postinstall.sh
    rm -f /root/postinstall.sh
    systemctl disable rc-local
    shutdown -r now
fi
exit 0
